# nestshare.service
[Unit]
Description=NestShare Dashboard
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/nestshare
ExecStart=/opt/nestshare/venv/bin/python app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
